# non-regression-tests/crontab
* 4 * * * /usr/local/bin/non-regression-tests 2>&1 > /var/log/non-regression-tests/$(date -Iseconds)
